# installation/supershy-daemon.service
[Unit]
Description=Supershy

[Service]
ExecStart=supershy
Restart=always

[Install]
WantedBy=default.target